# docker/base/etc/services.d/node/run
#!/usr/bin/with-contenv sh
cd /var/www

/usr/local/bin/npm run dev
